# crontab
*/1 * * * * /usr/local/bin/python /app/importer.py >> /var/log/cron.log 2>&1
